# cron/crontab
# START CRON JOB
@reboot /cron/cronjob.sh > /proc/1/fd/1 2>/proc/1/fd/2
# STOP CRON JOB
